# recipes-bsp/misc-config/files/update-misc-config.service
[Unit]
Description=update-misc-config

[Service]
ExecStart=/usr/bin/update-misc-config.sh

[Install]
WantedBy=multi-user.target
