# recipes-core/som-carrier-autoload/files/som-carrier-autoload.service
[Unit]
Description=som-carrier-autoload

[Service]
ExecStart=/usr/bin/som-carrier-autoload.sh

[Install]
WantedBy=multi-user.target
