# recipes-utils/som-dashboard/files/som-dashboard.service
[Unit]
Description=som-dashboard-init

[Service]
ExecStart=/usr/bin/som-dashboard-init

[Install]
WantedBy=multi-user.target
